# overlay/etc/services.d/heartbeat/run
#!/usr/bin/with-contenv sh

exec /scripts/heartbeat.sh